# engineering/tools/vps-deploy/systemd/iris-weekly-review.service
[Unit]
Description=Iris Weekly Review
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=eric
WorkingDirectory=/home/eric/OphidianAI
ExecStart=/opt/iris/scripts/weekly-review.sh
TimeoutStartSec=600
