# engineering/tools/vps-deploy/systemd/iris-follow-up-check.timer
[Unit]
Description=Iris Follow-Up Check -- Daily 9:00 AM ET

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
