# engineering/tools/vps-deploy/systemd/iris-follow-up-check.service
[Unit]
Description=Iris Follow-Up Check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=eric
WorkingDirectory=/home/eric/OphidianAI
ExecStart=/opt/iris/scripts/follow-up-check.sh
TimeoutStartSec=300
